# config/nginx.conf
# file: /etc/init/nginx.conf
description "NGiNX starter"
start on (local-filesystems and runlevel [2345])
stop on runlevel [016]
expect fork
respawn
exec /usr/local/sbin/nginx
